# packaging/systemd/some-ting-listen.service
# some-ting listener — systemd *user* service (Linux/X11).
# Install:
#   cargo install --path listener        # -> ~/.cargo/bin/some-ting-listen
#   mkdir -p ~/.config/systemd/user
#   cp packaging/systemd/some-ting-listen.service ~/.config/systemd/user/
#   systemctl --user import-environment DISPLAY XAUTHORITY   # X11 access for focus/inject
#   systemctl --user enable --now some-ting-listen
#   journalctl --user -u some-ting-listen -f                 # logs
#
# If the default input isn't the TING front-mic, add: --device <name-substring>
# (find it with: some-ting-listen --list-devices)

[Unit]
Description=some-ting: TING Quindar tones -> Claude voice push-to-talk
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.cargo/bin/some-ting-listen
Restart=always
RestartSec=2
# X11 focus/injection needs the session display; imported above, or set here:
# Environment=DISPLAY=:0
# Environment=XAUTHORITY=%h/.Xauthority

[Install]
WantedBy=graphical-session.target
